# deployments/resources/pager-forwarding.service
[Unit]
Description=Pager Relais Forwarding Service

[Service]
User=pager_forwarding
Group=pager_forwarding
ExecStart=/usr/pager-forwarding/main